# prowlarr/prowlarr.service
[Unit]
Description=Prowlarrr
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker compose -f /home/gart/prowlarr/docker-compose.yml up
ExecStop=/usr/bin/docker compose -f /home/gart/prowlarr/docker-compose.yml down

[Install]
WantedBy=default.target
