# bazarr/bazarr.service
[Unit]
Description=Bazarr
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker compose -f /home/gart/bazarr/docker-compose.yml up
ExecStop=/usr/bin/docker compose -f /home/gart/bazarr/docker-compose.yml down

[Install]
WantedBy=default.target
